# scripts/nanocode.service
# /etc/systemd/system/nanocode.service
#
# System-mode router. Workers are spawned per-user by `nanocode login`
# (via the setuid helper) and are NOT children of this unit.
#
# Pre-install requirements (handled by scripts/install.sh):
#   - group `nanocode` exists (service account only — NOT a login group;
#     users do NOT need to be members)
#   - user `nanocode` exists, GID=nanocode, /usr/sbin/nologin
#   - /usr/lib/nanocode/ contains server/, worker/, public/, vendor/,
#     and nanocode-spawn (mode 4755 root:root)
#   - /usr/local/bin/nanocode is the CLI (mode 0755, NOT setuid)
#
# Per-user runtime directories (/run/nanocode/u<uid>/) are created by
# the setuid helper at first `nanocode login` invocation. They are
# user-owned with group=nanocode + SGID, so worker sockets inside
# inherit the nanocode group automatically — the router can connect
# without any user needing to be a member of the nanocode group.

[Unit]
Description=Nanocode multi-user dispatcher
Documentation=file:///usr/share/doc/nanocode/system-mode-design.md
After=network.target

[Service]
Type=simple
User=nanocode
Group=nanocode
RuntimeDirectory=nanocode
# 0755 lets any user traverse into /run/nanocode/ to reach their own
# /run/nanocode/u<uid>/ subdir (created by the setuid helper). The
# per-uid subdirs are mode 02750 — only the owner + the nanocode group
# (the router) can traverse into them.
RuntimeDirectoryMode=0755
# Preserve /run/nanocode/ across `systemctl restart` so per-user worker
# sockets (u<uid>/worker.sock) survive a router restart — the workers
# stay alive and reattach via auto-reconnect. On full shutdown the
# directory is cleared as normal.
RuntimeDirectoryPreserve=restart

# /var/lib/nanocode/ holds persisted session tokens (3-day rolling TTL)
# so users stay logged in across router restarts. Mode 0700 keeps the
# file inaccessible to anyone but the router itself.
StateDirectory=nanocode
StateDirectoryMode=0700

ExecStart=/usr/bin/node /usr/lib/nanocode/server/index.js
Environment=NANOCODE_SYSTEM=1
Environment=HOST=0.0.0.0
Environment=PORT=2333
WorkingDirectory=/usr/lib/nanocode

Restart=on-failure
RestartSec=2s

# --- Hardening ---
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
ProtectProc=invisible
LockPersonality=yes
# MemoryDenyWriteExecute=yes blocks V8's JIT mappings on Node <20 (causes
# "Fatal javascript OOM" at boot). Re-enable when the deployed node is
# 20.6+ AND uses PKEY-aware JIT.
# MemoryDenyWriteExecute=yes
RestrictRealtime=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=yes
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
CapabilityBoundingSet=
AmbientCapabilities=
UMask=0077

# The router publishes /run/nanocode/router.sock; workers create their
# own sockets at /run/nanocode/u-<user>.sock. The RuntimeDirectoryMode
# above (0775) lets nanocode-group members create files there.

[Install]
WantedBy=multi-user.target
